# claw.service
[Unit]
Description=Claw AI Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
ExecStart=/home/YOUR_USERNAME/.local/bin/claw gateway
Restart=on-failure
RestartSec=5
KillMode=control-group
TimeoutStopSec=30
Environment=PATH=/home/YOUR_USERNAME/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
