# claw.service
[Unit]
Description=Claw AI Agent (gateway + WebUI on a single port)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
# The merged claw binary launches the gateway, the WebUI HTTP layer, and the
# session API on the same port (cfg.Gateway.Port, default 18790). No
# subcommand is needed — `claw` with no arguments is the supported entry
# point. The legacy `claw gateway` alias is still accepted.
ExecStart=/home/YOUR_USERNAME/.local/bin/claw
Restart=on-failure
RestartSec=5
KillMode=control-group
TimeoutStopSec=30
Environment=PATH=/home/YOUR_USERNAME/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
